# cron/crontab
* * * * * /bin/sh /opt/cron/limpar-processos.sh > /proc/1/fd/1 2>/proc/1/fd/2
