# cron/crontab
* * * * * root sh /opt/cron/limpar-processos.sh 2>&1
